# etc/crontabs/root
*/10 * * * * get_sysinfo_json.sh status online
